# src/openrc/agnikai
#!/sbin/openrc-run

# copy this file to /etc/init.d/agnikai
# copy tty.colors to /tty.colors

depend() {
  after modules
}

start() {
  ebegin "Starting agni-kai"
  start-stop-daemon --start --exec /usr/bin/setvtrgb /tty.colors
  eend $?
}

begin() {
  ebegin "Stopping agni-kai"
  start-stop-daemon --stop --exec /usr/bin/setvtrgb /tty.colors
  eend $?
}
